# debian/powerwaked.service
[Unit]
Description=PowerWake Server

[Service]
Type=forking
ExecStart=/usr/sbin/powerwaked

[Install]
WantedBy=multi-user.target
