# debian/powernap.upstart
# powernap

description "PowerNap"
author "Dustin Kirkland <kirkland@canonical.com>"

start on runlevel [2345]

expect daemon

pre-start exec /usr/share/powernap/powernap-ethtool

pre-stop script
	powerwake-now || true
end script

exec /usr/sbin/powernapd
